# etc/systemd/system/cc-archiver.service
[Unit]
Description=CCArchiver Services

[Service]
Type=oneshot
ExecStart=/bin/true
RemainAfterExit=yes

[Install]
WantedBy=cc-archiver.target
